# conf/deb/postinst
systemctl enable stream.service
update-rc.d stream defaults